# recipes-graphics/benchmark-demo/benchmark-demo-init/rz_benchmark_demo.service
[Unit]
Description=RZ Linux Benchmark Demo

[Service]
Type=simple
ExecStart=/home/root/start_demo.sh

[Install]
WantedBy=multi-user.target
